# src/etc/systemd/system/timekpr-checker.service
[Unit]
Requires=timekpr.service
After=suspend.target

[Service]
ExecStart=echo

[Install]
WantedBy=suspend.target
